# 4.5.10-php7.1/scripts/supervisord-zwordpress.conf
[program:wordpress-setup]
command=sh -x /setup-wordpress.sh
startsecs=1
autorestart=false
startretries=1